# deploy/com.jaywest.multi-agent-observability.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jaywest.multi-agent-observability</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jaywest/.bun/bin/bun</string>
        <string>run</string>
        <string>src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jaywest/Multi-Agent-Observability/claude-code-hooks-multi-agent-observability/apps/server</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jaywest/Library/Logs/multi-agent-observability.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jaywest/Library/Logs/multi-agent-observability-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jaywest/.bun/bin:/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jaywest</string>
    </dict>
</dict>
</plist>
